# client/daemon/python3/systemd/minute-timer-polousers.target
[Unit]
Description=Minute Timer Target for Polousers
StopWhenUnneeded=yes