# client/daemon/python3/systemd/minute-timer-polousers.timer
[Unit]
Description=Minute Timer

[Timer]
OnBootSec=1min
Unit=minute-timer-polousers.target

[Install]
WantedBy=basic.target